# gatekeeper_init
#!/bin/sh /etc/rc.common
#
# The MIT License (MIT)
# Copyright (c) 2026 Naresh Mehta (https://www.naresh.se/)
#
# Permission is hereby granted, free of charge, to any person obtaining a copy
# of this software and associated documentation files (the "Software"), to deal
# in the Software without restriction, including without limitation the rights
# to use, copy, modify, merge, publish, distribute, sublicense, and/or sell
# copies of the Software, and to permit persons to whom the Software is
# furnished to do so, subject to the following conditions:
#
# The above copyright notice and this permission notice shall be included in
# all copies or substantial portions of the Software.
#
# THE SOFTWARE IS PROVIDED "AS IS", WITHOUT WARRANTY OF ANY KIND, EXPRESS OR
# IMPLIED, INCLUDING BUT NOT LIMITED TO THE WARRANTIES OF MERCHANTABILITY,
# FITNESS FOR A PARTICULAR PURPOSE AND NONINFRINGEMENT. IN NO EVENT SHALL THE
# AUTHORS OR COPYRIGHT HOLDERS BE LIABLE FOR ANY CLAIM, DAMAGES OR OTHER
# LIABILITY, WHETHER IN AN ACTION OF CONTRACT, TORT OR OTHERWISE, ARISING FROM,
# OUT OF OR IN CONNECTION WITH THE SOFTWARE OR THE USE OR OTHER DEALINGS IN THE
# SOFTWARE.
#

# gatekeeper_init - OpenWrt init script for gatekeeper MAC synchronization
#
# This script runs at system startup (START=95) to initialize the gatekeeper
# firewall rules by populating nftables sets with MAC addresses from UCI configuration.
# Syncs both static DHCP leases and blacklist MACs.
#
# Key Features:
# - Runs during late boot phase to ensure network services are ready
# - Clears existing sets before repopulating (prevents stale entries)
# - Parses UCI DHCP configuration to extract static MAC assignments
# - Parses UCI gatekeeper configuration to extract blacklist MACs
# - Adds each MAC to the appropriate nftables firewall set for access control
#
# Integration Points:
# - Works with OpenWrt's procd init system
# - Modifies nftables firewall rules via 'fw4' table in 'inet' family
# - Reads from UCI 'dhcp' configuration namespace
#
# Security Considerations:
# - Only processes legitimate static DHCP leases from system config
# - No external input validation needed as UCI config is trusted
# - Failures during MAC addition are logged by nftables but don't stop the script
#
# Extensibility Notes:
# - To add MACs from additional sources, extend the STATIC_MACS extraction logic
# - For custom timing, adjust START parameter (higher values = later execution)
# - To support IPv6, add similar logic for static IPv6 assignments

# OpenWrt init system parameters
START=95  # Late boot position (90-99 range for network-related services)

# Main service start function called by procd
start() {
    # === SYNC STATIC DHCP MACS ===
    # Step 1: Clear the static MAC set to prevent accumulation of stale entries
    nft flush set inet fw4 static_macs

    # Step 2: Extract static MAC addresses from UCI DHCP configuration
    STATIC_MACS=$(uci show dhcp | grep ".mac=" | awk -F"='" '{print $2}' | tr -d "'")

    # Step 3: Iterate through each static MAC and add to firewall set
    for mac in $STATIC_MACS; do
        nft add element inet fw4 static_macs { $mac }
    done

    logger -t gatekeeper_init "Static MACs synchronized"

    # === SYNC BLACKLIST MACS ===
    # Step 4: Clear the blacklist MAC set
    nft flush set inet fw4 blacklist_macs 2>/dev/null

    # Step 5: Extract blacklist MAC addresses from UCI gatekeeper configuration
    BLACKLIST_MACS=$(uci show gatekeeper.blacklist 2>/dev/null | grep "\.mac=" | cut -d"'" -f2)

    # Step 6: Iterate through each blacklist MAC and add to firewall set
    for mac in $BLACKLIST_MACS; do
        [ -z "$mac" ] && continue
        nft add element inet fw4 blacklist_macs { $mac } 2>/dev/null
    done

    logger -t gatekeeper_init "Blacklist MACs synchronized"
}
